# kadeboot-firstrun.service
[Unit]
Description=KADEBOOT First Run Setup - Opens KAES-ARCH Configuration in Konsole
After=graphical-session.target
Wants=graphical-session.target
ConditionPathExists=!/home/%i/.config/kadeboot-setup-complete
ConditionPathExists=/home/%i/KAES-ARCH/post.sh

[Service]
Type=oneshot
User=%i
Group=users
Environment=DISPLAY=:0
Environment=HOME=/home/%i
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStartPre=/bin/sleep 3
ExecStart=/usr/bin/konsole --hold -e /home/%i/KAES-ARCH/post.sh
ExecStartPost=/usr/bin/touch /home/%i/.config/kadeboot-setup-complete
RemainAfterExit=no

[Install]
WantedBy=graphical-session.target